# data/downloader.cron
57 06 * * * /app/downloader_launcher.sh
